# root/etc/systemd/system/clamd@.service.d/reload.conf
[Service]
ExecReload=/usr/bin/clamdscan --config-file=/etc/clamd.d/%i.conf --reload
