# root/etc/cron.hourly/freshclam-random
#!/bin/bash
/bin/sleep $[ $RANDOM % 5 ] && exec /etc/cron.daily/freshclam
